# scripts/s6/act_runner/run
#!/usr/bin/env bash

s6-svwait -U /etc/s6/docker

exec run.sh
