# scripts/s6/docker/run
#!/usr/bin/env bash

exec s6-notifyoncheck dockerd-entrypoint.sh
