# scripts/s6/docker/run
#!/usr/bin/env bash

exec dockerd-entrypoint.sh
